# rootfs/etc/services.d/clickhouse/run
#!/command/with-contenv bash
set -euo pipefail

# shellcheck source=/dev/null
. /opt/signoz-aio/lib/env.sh
normalize_blank_env "${AIO_BLANK_AS_UNSET_VARS[@]}"

if is_true "${SIGNOZ_USE_EXTERNAL_CLICKHOUSE:-false}"; then
    aio_log "External ClickHouse mode enabled; skipping bundled ClickHouse."
    exec sleep infinity
fi

ulimit -n 262144 || true
ulimit -u 65535 || true

exec su -s /bin/bash clickhouse -c 'exec clickhouse-server --config-file=/etc/clickhouse-server/config.xml'
