# rootfs/etc/services.d/otel-collector/run
#!/command/with-contenv bash
set -euo pipefail

MIGRATION_MARKER="/appdata/.telemetrystore-migrations-complete"

wait_for_clickhouse() {
    until curl -fsS http://127.0.0.1:8123/ping >/dev/null 2>&1; do
        sleep 2
    done
}

wait_for_signoz() {
    until curl -fsS http://127.0.0.1:8080/api/v1/health >/dev/null 2>&1; do
        sleep 2
    done
}

COLLECTOR_BIN="/opt/signoz-otel-collector/signoz-otel-collector"
COLLECTOR_CONFIG="/opt/signoz-aio/config/otel-collector-config.yaml"

CLICKHOUSE_DSN="${SIGNOZ_OTEL_COLLECTOR_CLICKHOUSE_DSN:-tcp://127.0.0.1:9000}"
CLICKHOUSE_CLUSTER="${SIGNOZ_OTEL_COLLECTOR_CLICKHOUSE_CLUSTER:-cluster}"
CLICKHOUSE_REPLICATION="${SIGNOZ_OTEL_COLLECTOR_CLICKHOUSE_REPLICATION:-false}"

export SIGNOZ_OTEL_COLLECTOR_CLICKHOUSE_DSN="${CLICKHOUSE_DSN}"
export SIGNOZ_OTEL_COLLECTOR_CLICKHOUSE_CLUSTER="${CLICKHOUSE_CLUSTER}"
export SIGNOZ_OTEL_COLLECTOR_CLICKHOUSE_REPLICATION="${CLICKHOUSE_REPLICATION}"

wait_for_clickhouse

if [ ! -f "${MIGRATION_MARKER}" ]; then
    "${COLLECTOR_BIN}" migrate bootstrap
    "${COLLECTOR_BIN}" migrate sync up
    "${COLLECTOR_BIN}" migrate async up
    touch "${MIGRATION_MARKER}"
fi

wait_for_signoz

"${COLLECTOR_BIN}" migrate sync check

exec "${COLLECTOR_BIN}" \
  --config "${COLLECTOR_CONFIG}"
